# blfs-servicescripts/services/proftpd/run
#!/bin/sh
exec 2>&1
exec proftpd -n -d 1