# blfs-servicescripts/services/winbindd/run
#!/bin/sh

sv check smbd nmbd >/dev/null || exit 0

exec /usr/sbin/winbindd