# blfs-servicescripts/services/svn/run
#!/bin/sh
exec /usr/bin/svnserve -d --foreground 2>&1