# blfs-servicescripts/services/postfix/run
#!/bin/sh
/usr/sbin/postfix check || exit 1
exec /usr/lib/postfix/master
