# blfs-servicescripts/services/postgresql/run
#!/bin/sh
exec setuidgid postgres /usr/lib/postgresql/bin/postmaster \
  -D /var/lib/postgres/data 2>&1