# blfs-servicescripts/services/qpopper/run
#!/bin/sh
exec 2>&1
exec tcpserver -v 0 pop3 popper -f /etc/qpopper.config -l0 -s -d
