# blfs-servicescripts/services/httpd/finish
#!/bin/sh
exec /usr/sbin/apachectl -k stop