# blfs-servicescripts/services/httpd/run
#!/bin/sh
exec 2>&1
mkdir -p /var/run/httpd||true
exec httpd -eDEBUG -DNO_DETACH -DFOREGROUND