# lfs-servicescripts/services/agetty-tty5/run
#!/bin/sh
exec chpst -P /sbin/agetty tty5 linux