# lfs-servicescripts/services/agetty-tty4/run
#!/bin/sh
exec chpst -P /sbin/agetty tty4 linux