# lfs-servicescripts/services/agetty-tty3/run
#!/bin/sh
exec chpst -P /sbin/agetty tty3 linux