# blfs-servicescripts/services/autofs/run
#!/bin/sh
[ -f /etc/sysconfig/autofs ] && . /etc/sysconfig/autofs
[ -n "$TIMEOUT" ] && daemonoptions="--timeout=$TIMEOUT $daemonoptions"
exec /usr/sbin/automount -f $daemonoptions >/dev/null 2>&1