# blfs-servicescripts/services/krb5/run
#!/bin/sh
exec /usr/sbin/krb5kdc
exec /usr/sbin/kadmind
exec /usr/sbin/kpropd -S