# blfs-servicescripts/services/smartd/run
#!/bin/sh
exec 2>&1
exec smartd -n -p /run/smartd.pid -q never 2>&1