# blfs-servicescripts/services/ntpd/run
#!/bin/sh
exec 2>&1
exec ntpd --nofork -g -u ntp:ntp