# blfs-servicescripts/services/netfs/run
#!/bin/sh
exec /bin/mount -a -O _netdev
exec /usr/bin/pause