# lfs-servicescripts/services/agetty-tty2/run
#!/bin/sh
exec chpst -P /sbin/agetty tty2 linux