# blfs-servicescripts/services/sysstat/run
#!/bin/sh

exec /usr/lib/sa/sadc -F -L -

exec /usr/bin/pause