# lfs-servicescripts/services/agetty-tty6/run
#!/bin/sh
exec chpst -P /sbin/agetty tty6 linux