# blfs-servicescripts/services/soprano/run
#!/bin/sh
[ -r /etc/sysconfig/soprano ] && . /etc/sysconfig/soprano
HOME=$SOPRANO_STORAGE
exec /usr/bin/sopranod \
	--storagedir $SOPRANO_STORAGE \
	--backend $SOPRANO_BACKEND \
	$SOPRANO_OPTIONS