# blfs-servicescripts/services/acpid/run
#!/bin/sh
exec 2>&1
exec acpid -c /etc/acpi/events -f