# lfs-servicescripts/services/agetty-tty1/run
#!/bin/sh
exec chpst -P /sbin/agetty --noclear 38400 tty1 linux