# blfs-servicescripts/services/iptables/finish
#!/bin/bash

iptables -F
iptables -P INPUT DROP
iptables -P FORWARD DROP
iptables -P OUTPUT DROP