# blfs-servicescripts/services/php-fpm/run
#!/bin/sh

exec /usr/sbin/php-fpm --daemonize --fpm-config /etc/php-fpm.conf --pid /var/run/php-fpm.pid