# blfs-servicescripts/services/dhcpd/run
#!/bin/sh
exec 2>&1
exec dhcpd -f -d -cf ./config eth0