# blfs-servicescripts/services/sendmail/run
#!/bin/sh

sv check slapd postgresql mysql >/dev/null || exit 0

exec /usr/sbin/sendmail -bs -bd -q5m start