# blfs-servicescripts/services/nfs-client/run
#!/bin/sh
sv check rpcbind >/dev/null || exit 1
exec /usr/sbin/rpc.statd